# temp_mon.cron
*/5 * * * * root /usr/local/bin/temp_mon.sh
